# modules/obnam/etc-template/cron.d/obnam-run-backups
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games

0 3 * * * root @@SCRIPT_DIR@@/sanei module obnam backup-all
