# modules/obnam/etc-template/cron.d/obnam-run-backups
0 3 * * * root sanei module obnam run-backup