# modules/obnam/etc-template/cron.d/obnam-run-backups
0 3 * * * root bash @@MODULES_DIR@@/run-backup.sh