# modules/lxc-common/etc-template/init/lxc-custom.conf
start on mounted MOUNTPOINT=/run

script
    ln -s /shared/run/mysqld /run/mysqld
    ln -s /shared/run/postgresql /run/postgresql
end script