# modules/mysql-xtrabackup/etc-template/cron.d/mysql-xtrabackup
0 1 * * * root @@SCRIPT_DIR@@/sanei module mysql-xtrabackup backup