# modules/scheduled-tasks/etc-template/cron.d/sanei-scheduled-tasks
0 6 * * * root cd @@SCRIPT_DIR@@; run-parts @@MODULE_DIR@@/tasks