# modules/website-archiver/etc-template/cron.d/website-archiver-delicious
0 * * * * @@DELICIOUS_SYSTEM_USER@@ bash @@MODULE_DIR@@/delicious.sh
