# etc/init/lxc-custom.conf
start on mounted MOUNTPOINT=/run

exec ln -s /shared/run/mysqld /run/mysqld

